# packages/debian/DEBIAN/prerm
#!/bin/bash

set -e

# Stop and disable systemd service if it's running
if systemctl is-active --quiet apk-reverse-tool.service; then
    systemctl stop apk-reverse-tool.service
    systemctl disable apk-reverse-tool.service
fi

exit 0